# deploy/base/systemd/geomark-base.service
[Unit]
Description=GeoMark Base Station
After=network.target dev-ttyAMA0.device
Wants=dev-ttyAMA0.device

[Service]
Type=simple
User=pi
ExecStart=/usr/local/bin/geomark --mode base --config /etc/geomark/geomark.conf
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target