# systemd/geomark-base.service
[Unit]
Description=GeoMark Base Station
After=network.target

[Service]
Type=simple
User=pi
ExecStart=/home/pi/geomark/scripts/start_base.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target