# systemd/geomark-base.service
[Unit]
Description=GeoMark Base Station
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/geomark --mode base --config /etc/geomark/geomark.conf
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target