# deploy/handheld/systemd/geomark-connect.service
[Unit]
Description=GeoMark — connect handheld to rover AP
After=network.target NetworkManager.service NetworkManager-wait-online.service
Wants=NetworkManager.service
Requires=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/geomark-connect.sh
TimeoutStartSec=90

[Install]
WantedBy=multi-user.target